# supervisord.conf

[supervisord]
logfile=/var/log/supervisord.log ; (main log file;default $CWD/supervisord.log)
pidfile=/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=true                ; (start in foreground if true;default false)
childlogdir=/var/log/supervisor ; ('AUTO' child log dir, default $TEMP)

[program:avahi]
command=avahi-daemon --no-drop-root
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:snapserver]
command=snapserver -s "spotify:///usr/local/bin/librespot?name=Spotify&username=%(ENV_SPOTIFYUSER)s&password=%(ENV_SPOTIFYPASS)s" -s "airplay:///usr/bin/shairport-sync?name=Airplay"
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

